# scripts/launchd/com.postcabinets.linear-agent-loop.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.postcabinets.linear-agent-loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/apple/claude for me/scripts/linear-agent-loop.mjs</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/apple/Library/Logs/linear-agent-loop.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/apple/Library/Logs/linear-agent-loop.error.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/apple</string>
  </dict>
</dict>
</plist>
