# scripts/launchd/com.postcabinets.nightly-scan.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.postcabinets.nightly-scan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/apple/claude for me/scripts/notion-nightly-scan.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/apple/claude for me</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/apple/Library/Logs/postcabinets/nightly-scan.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/apple/Library/Logs/postcabinets/nightly-scan.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/apple</string>
  </dict>
</dict>
</plist>
