# ops/scripts/launchd/com.postcabinets.discord-bot-daemon.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.postcabinets.discord-bot-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nobu/.nvm/versions/node/v20.19.5/bin/node</string>
    <string>/Users/nobu/claude for me/scripts/discord-bot-daemon.mjs</string>
  </array>

  <!-- クラッシュしたら自動再起動 -->
  <key>KeepAlive</key>
  <true/>

  <!-- ログイン時に自動起動 -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/nobu/Library/Logs/discord-bot-daemon.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nobu/Library/Logs/discord-bot-daemon.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/nobu/.nvm/versions/node/v20.19.5/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/apple</string>
  </dict>
</dict>
</plist>
