# ops/scripts/launchd/com.postcabinets.discord-morning.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.postcabinets.discord-morning</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nobu/.nvm/versions/node/v20.19.5/bin/node</string>
    <string>/Users/nobu/claude for me/scripts/discord-morning.mjs</string>
  </array>

  <!-- 毎朝 5:00 に実行 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/nobu/Library/Logs/discord-morning.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nobu/Library/Logs/discord-morning.error.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/nobu/.nvm/versions/node/v20.19.5/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/apple</string>
  </dict>
</dict>
</plist>
